# backend/scraper/crontab
0 3 * * * root cd /app && python -m scraper.run >> /var/log/scraper.log 2>&1
